# deploy/usipipo-backend.service
[Unit]
Description=uSipipo Backend API Service
Documentation=https://github.com/uSipipo-Team/usipipo-backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
User=mowgli
Group=mowgli
WorkingDirectory=/home/mowgli/usipipo/usipipo-backend

# Environment
Environment="PATH=/home/mowgli/usipipo/usipipo-backend/.venv/bin"
EnvironmentFile=-/home/mowgli/usipipo/usipipo-backend/.env

# Execution - Uses API_HOST and API_PORT from .env
ExecStart=/home/mowgli/usipipo/usipipo-backend/.venv/bin/uvicorn src.main:app --host ${API_HOST:-0.0.0.0} --port ${API_PORT:-8000} --workers 4
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/mowgli/usipipo/usipipo-backend/logs /home/mowgli/usipipo/usipipo-backend/temp /home/mowgli/usipipo/usipipo-backend/static

# Resource limits
LimitNOFILE=65535
Nice=-5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=usipipo-backend

[Install]
WantedBy=multi-user.target
